# udev_rules/99-keychron_k2_bluetooth.rules
ACTION=="add", ATTRS{uniq}=="dc:2c:26:e7:bf:5f", RUN+="/usr/local/bin/udev_keyboard_hook"
